# supervisord.conf
[supervisord]
nodaemon=true

[program:custom1]
command=/run.sh
user=grafana

#[program:Argo]
#command=/usr/share/argo-tunnel.sh
